# systemd/pictl-shim.service
[Unit]
Description=pictl localhost HTTP shim
Documentation=https://github.com/simon-lehmann/pictl
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pictl
ExecStart=/usr/bin/python3 /home/pi/pictl/shim.py
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Env defaults. Override with `systemctl edit pictl-shim` if needed.
Environment=PICTL_SHIM_HOST=127.0.0.1
Environment=PICTL_SHIM_PORT=8080
Environment=PICTL_SHIM_TIMEOUT=10

# Hardening. The shim only needs to read ~/.pictl, spawn pictl.py, and
# bind 127.0.0.1:8080. Everything else can be locked down.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/.pictl
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
